# simple-webcam.service
[Unit]
Description=Simple Webcam

[Service]
ExecStart=/home/pi/simple-webcam/simple-webcam

[Install]
WantedBy=multi-user.target